# deploy/systemd/mk04-source-input.service
[Unit]
Description=mk0.4 source-input ingestion service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mk04
Group=mk04
WorkingDirectory=/opt/mk04/prod/current
Environment=MK04_ENV=prod
EnvironmentFile=/etc/mk04/prod/env
EnvironmentFile=-/etc/mk04/prod/services/source-input.env
ExecStart=/opt/mk04/prod/current/deploy/scripts/run-input-service.sh prod
Restart=always
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mk04-source-input

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

# yt-dlp / download paths may open many handles.
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
