# deploy/cron/mk04.crontab
# mk0.4 autonomous mode cron table (Reliability & Recovery Phase 9)
#
# Install (survives reboot; cron starts with the system):
#   sudo -u mk04 crontab /opt/mk04/prod/current/deploy/cron/mk04.crontab
# Or system-wide drop-in:
#   sudo cp /opt/mk04/prod/current/deploy/cron/mk04.cron.d /etc/cron.d/mk04
#
# Pipeline scheduling is intentionally thin:
#   cron → scripts/ops/run-scheduled.sh → scripts/ops/run-pipeline.sh
#
# run-pipeline enforces config validation, boot readiness, execution lock,
# run records, and runtime scheduler enable/disable. Do not call /run-funnel
# or Python pipeline modules directly from cron.
#
# Times are UTC unless the host crond is configured otherwise.
# One pipeline line per active source-input funnel_id (stagger if multiple).

MK04_ENV=prod

# --- Scheduled pipeline (shared entrypoint) ----------------------------------
0 8 * * * /opt/mk04/prod/current/scripts/ops/run-scheduled.sh prod mfm_business_ai_001 >/dev/null 2>&1

# --- Handoff retry sweeper ---------------------------------------------------
# video-automation → output-funnel registration retry (not pipeline scheduling).
*/10 * * * * MK04_ENV=prod PIPELINE_CONFIG_PATH=/etc/mk04/prod/video-automation/pipeline_config.json OUTPUT_FUNNEL_URL=http://127.0.0.1:5055 VIDEO_AUTOMATION_JOBS_DIR=/var/lib/mk04/prod/video-automation/jobs /usr/bin/python3 /opt/mk04/prod/current/video-automation/scripts/handoff_sweeper.py --quiet >/dev/null 2>&1

# --- Database backup (Storage Phase 10) --------------------------------------
# SQLite snapshot only. Expired backups use retention (database_backups_days).
0 3 * * * /opt/mk04/prod/current/scripts/ops/run-database-backup.sh prod >/dev/null 2>&1

# --- Log rotation (Storage Phase 9) ------------------------------------------
# Bounds active project logs. Expired rotated logs use retention (logs_days).
15 3 * * * /opt/mk04/prod/current/scripts/ops/run-log-rotation.sh prod >/dev/null 2>&1

# --- Scheduled retention (Storage Phase 8) -----------------------------------
# Thin trigger only. Mode/enablement come from storage.schedule in config.
# Production defaults to dry_run. Does not call the legacy retention-sweeper.
30 3 * * * /opt/mk04/prod/current/scripts/ops/run-scheduled-retention.sh prod >/dev/null 2>&1

# --- Scheduled retention (DEV, optional) -------------------------------------
# 30 4 * * * /home/maguireltd/mk1-automated-clipping/scripts/ops/run-scheduled-retention.sh dev >/dev/null 2>&1

# --- Health watchdog ---------------------------------------------------------
*/15 * * * * /opt/mk04/prod/current/deploy/scripts/watchdog.sh prod >/dev/null 2>&1
